# salt/etc/supervisor/conf.d/salt-master.conf
[program:salt-master]
command=sh /srv/runinenv /srv/virtualenv/ salt-master
process_name=salt-master
user=root
