# salt/etc/supervisor/conf.d/salt-minion.conf
[program:salt-minion]
command=sh /srv/runinenv /srv/virtualenv/ salt-minion
process_name=salt-minion
user=root
